# src/assets/slidewsl/staging/context/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:dockerd]
command=/usr/bin/dockerd
autostart=true
autorestart=true
